# private_dot_config/systemd/user/trash-empty.timer
[Unit]
Description=Trash empty daily timer

[Timer]
OnCalendar=weekly
Persistent=true
Unit=trash-empty.service

[Install]
WantedBy=timers.target

