# others/home-dauliac-nextcloud.automount
[Unit]
Description=Mount WebDAV Service
After=network-online.target
Wants=network-online.target

[Automount]
Where=/home/dauliac/nextcloud
TimeoutIdleSec=300

[Install]
WantedBy=remote-fs.target

